# deploy/verify-valyd-scheduler.conf
; /etc/supervisor/conf.d/verify-valyd-scheduler.conf
; Runs Laravel's scheduler, which fires `verify:expire-stale` every minute to
; expire timed-out sessions and emit their webhooks.
;
; (Alternative: a system cron entry —
;   * * * * * cd /var/www/pollus_main_servers/verify.pollus.tech/backend && php artisan schedule:run >> /dev/null 2>&1
; )
;
; Install:
;   sudo cp deploy/verify-valyd-scheduler.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread && sudo supervisorctl update
;   sudo supervisorctl start verify-valyd-scheduler

[program:verify-valyd-scheduler]
command=php /var/www/pollus_main_servers/verify.pollus.tech/backend/artisan schedule:work
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/pollus_main_servers/verify.pollus.tech/backend/storage/logs/scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
