# deploy/verify-valyd-worker.conf
; /etc/supervisor/conf.d/verify-valyd-worker.conf
; Processes the `database` queue — delivers signed verification webhooks
; (SendVerificationWebhookJob, with retries/backoff).
;
; Install:
;   sudo cp deploy/verify-valyd-worker.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread && sudo supervisorctl update
;   sudo supervisorctl start verify-valyd-worker:*

[program:verify-valyd-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/pollus_main_servers/verify.pollus.tech/backend/artisan queue:work database --queue=default --sleep=3 --tries=10 --max-time=3600 --timeout=120 --memory=256
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/pollus_main_servers/verify.pollus.tech/backend/storage/logs/worker.log
stopwaitsecs=3600
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
